# fibservice/supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=/bin/bash -c "python /usr/local/src/fibservice/fibService.py"
