# servicelocator/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:etcd]
command=/bin/bash -c "/usr/local/src/servicelocator/startEtcd.sh"
