# webapp/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:flask]
command=/bin/bash -c "/usr/local/src/webapp/startFlask.sh"
