# webservice/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug


[program:flask]
priority=990
command=/bin/bash -c "/usr/local/src/webservice/startFlask.sh"
